# luci-app-natter2/root/etc/uci-defaults/luci-natter2
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@natter2[-1]
	add ucitrack natter2
	set ucitrack.@natter2[-1].init=natter2
	commit ucitrack
EOF

exit 0
