# luci-mod-istorenext/root/etc/uci-defaults/50_luci-istorenext
#!/bin/sh

rm -f /tmp/luci-indexcache
/etc/init.d/nginx reload 2>/dev/null
exit 0
