# mmdvm-luci/luci-theme-bootstrap-mmdvm/root/etc/uci-defaults/30_luci-theme-bootstrap-mmdvm
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Bootstrap_MMDVM=/luci-static/bootstrap-mmdvm
	set luci.main.mediaurlbase=/luci-static/bootstrap-mmdvm
	commit luci
EOF
exit 0
