# vmease/files/vmease.init
#!/bin/sh /etc/rc.common

START=98
USE_PROCD=1

start_service() {
	procd_open_instance
	procd_set_param limits nofile="65535 65535"
	procd_set_param command /usr/sbin/vmeasedaemon
	procd_append_param command vmeasedaemon
	procd_set_param respawn
	procd_close_instance
}

service_triggers() {
	procd_add_reload_trigger "vmease"
}

