# mmdvm-luci/luci-app-mmdvm/root/etc/uci-defaults/40_luci-mmdvm
#!/bin/sh

sed -i "s/cgi-bin\/luci\/\"/cgi-bin\/luci\/mmdvm\/\"/g" /www/index.html
rm -rf /tmp/luci-*
exit 0