# phtunnel/files/etc/hotplug.d/iface/30-oray-phtunnel
#!/bin/sh

[ "$INTERFACE" == "wan" -a "$ACTION" == "ifup" ] && {
	logger -t oray 'wan interface up and phtunnel restart'
	/etc/init.d/phtunnel restart
}
