# luci-app-fakemesh/root/etc/uci-defaults/40_luci-app-fakemesh
touch /etc/config/fakemesh
touch /etc/config/wireless

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@fakemesh[-1]
	add ucitrack fakemesh
	set ucitrack.@fakemesh[-1].init=fakemesh
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
