# luci-app-gecoosac/root/etc/uci-defaults/luci-gecoosac
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@gecoosac[-1]
	add ucitrack gecoosac
	set ucitrack.@gecoosac[-1].init=gecoosac
	commit ucitrack
EOF

rm -rf /tmp/luci-*
exit 0
