# small-package/small-package/wiretrustee/files/wiretrustee.init
#!/bin/sh /etc/rc.common

USE_PROCD=1

START=99
STOP=15

NAME=wiretrustee

start_service() {
  procd_open_instance
  procd_set_param command /bin/sh -c "/usr/bin/$NAME service run"
  procd_set_param pidfile /var/run/$NAME.pid
  procd_close_instance
}
