# luci-app-openvpn-server/root/etc/uci-defaults/openvpn.rej
--- luci-app-openvpn-server/root/etc/uci-defaults/openvpn
+++ luci-app-openvpn-server/root/etc/uci-defaults/openvpn
@@ -1,5 +1,7 @@
 #!/bin/sh
 
+grep -qv "myvpn" /etc/config/openvpn && grep -q "myvpn" /etc/config/openvpn-opkg && cp -f /etc/config/openvpn-opkg /etc/config/openvpn
+
 openvpn_port="$(uci -q get openvpn.myvpn.port)"
 [ -z "$openvpn_port" ] && openvpn_port=1194
 
