# luci-app-homeproxy/root/etc/uci-defaults/luci-homeproxy-migration
#!/bin/sh

ucode "/etc/homeproxy/scripts/migrate_config.uc"

exit 0
