# my-default-settings/target/rockchip/armv8/files/etc/hotplug.d/iface/12-rk3328
#!/bin/sh

ETHTOOL=/usr/sbin/ethtool

case "$(board_name)" in
friendlyelec,nanopi-r2s)
	[ "$(uci get network.$INTERFACE.device)" = "eth1" ] && \
	[ "$ACTION" = "ifup" ] && \
	$ETHTOOL -K eth1 rx off tx off
	;;
friendlyelec,nanopi-r2c)
	[ "$(uci get network.$INTERFACE.device)" = "br-lan" ] && \
	[ "$ACTION" = "ifup" ] && \
	$ETHTOOL -K eth1 rx off tx off
	;;
esac
