# pgyvpn/files/etc/hotplug.d/iface/pgyvpn
#!/bin/sh

[ "$INTERFACE" == "wan" -a "$ACTION" == "ifup" ] && {
	/etc/init.d/pgyvpn restart
}
