# luci-app-routerdog/root/etc/uci-defaults/50_luci-routerdog
#!/bin/sh

uci set wizard.default.landing_page='routerdog'
uci commit wizard

rm -f /tmp/luci-indexcache
exit 0
