# luci-app-socat/root/etc/uci-defaults/luci-app-socat
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@socat[-1]
	add ucitrack socat
	set ucitrack.@socat[-1].init=luci_socat
	commit ucitrack
EOF

rm -rf /tmp/luci-*cache*
exit 0
