# luci-app-mosdns/root/etc/hotplug.d/iface/99-mosdns
#!/bin/sh
[[ "$ACTION" = ifup && "$(uci -q get mosdns.mosdns.enabled)" == 1 ]] && /etc/init.d/mosdns restart
