# luci-app-xray/core/root/etc/hotplug.d/iface/01-transparent-proxy-ipset
#!/bin/sh
logger -st transparent-proxy-ipset[$$] -p6 "$(ucode /usr/share/xray/default_gateway.uc)"
