# luci-app-xray/root/etc/hotplug.d/iface/01-transparent-proxy-ipset.fw3
#!/bin/sh
DEFAULT_GATEWAY=$(ip route get 1.1.1.1 | grep -Eoe "src (\d+.\d+.\d+.\d+)" | awk '{print $NF}')
if [ -n "$DEFAULT_GATEWAY" ] ; then
    logger -st transparent-proxy-ipset[$$] -p6 "default gateway available at $DEFAULT_GATEWAY"
    ipset -! restore <<-EOF
		create tp_spec_def_gw hash:net hashsize 64
		flush tp_spec_def_gw
		add tp_spec_def_gw $DEFAULT_GATEWAY
	EOF
else
    logger -st transparent-proxy-ipset[$$] -p6 "default gateway not available, please wait for interface ready"
fi
