# luci-app-xwan/files/root/etc/uci-defaults/40_luci-app-xwan

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@xwan[-1]
	add ucitrack xwan
	set ucitrack.@xwan[-1].init=xwan
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
