# small-package/luci-app-mail/root/etc/uci-defaults/4101-mail
#!/bin/sh
uci -q batch <<-EOF >/dev/null
	delete ucitrack.@mail[-1]
	add ucitrack mail
	set ucitrack.@mail[-1].init=mail
	commit ucitrack
EOF
exit 0