# small-package/small-package/luci-app-dnsmasq-ipset/files/root/etc/uci-defaults/luci-dnsmasq-ipset
#!/bin/sh
uci -q batch <<-EOF >/dev/null
    delete ucitrack.@dnsmasq-ipset[-1]
    add ucitrack dnsmasq-ipset
    set ucitrack.@dnsmasq-ipset[-1].init=dnsmasq-ipset
    commit ucitrack
EOF
exit 0
