# luci-app-wizard/root/etc/uci-defaults/96-wizard-defaults
NETMASK=$(uci -q get network.lan.netmask)
DNS=$(uci -q get network.lan.dns)
GATEWAY=$(uci -q get network.lan.gateway)
PROTO=$(uci -q get network.wan.proto)

uci -q set wizard.default=wizard
uci -q set wizard.default.lan_netmask="$NETMASK"
uci -q set wizard.default.lan_dns="$DNS"
uci -q set wizard.default.lan_gateway="$GATEWAY"
[ "$(uci -q get dhcp.lan.ignore)" ] && uci -q set wizard.default.dhcp="0"
uci -q set wizard.default.wan_proto="$PROTO"

uci commit wizard

exit 0
