# small-package/small-package/luci-app-aliyundrive-fuse/root/etc/uci-defaults/luci-aliyundrive-fuse
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@aliyundrive-fuse[-1]
	add ucitrack aliyundrive-fuse
	set ucitrack.@aliyundrive-fuse[-1].init=aliyundrive-fuse
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
