# diy/.packages/luci-app-wizard/root/etc/uci-defaults/96-wizard-defaults

SSID=$(uci get wireless.@wifi-iface[0].ssid 2>/dev/null)
KEY=$(uci get wireless.@wifi-iface[0].key 2>/dev/null)
IPADDR=$(uci get network.lan.ipaddr 2>/dev/null)
NETMASK=$(uci get network.lan.netmask 2>/dev/null)
PROTO=$(uci get network.wan.proto 2>/dev/null)
USERNAME=$(uci get network.wan.username 2>/dev/null)
PASSWORD=$(uci get network.wan.password 2>/dev/null)

uci set wizard.default=wizard
uci set wizard.default.wifi_ssid="$SSID"
uci set wizard.default.wifi_key="$KEY"
uci set wizard.default.lan_ipaddr="$IPADDR"
uci set wizard.default.lan_netmask="$NETMASK"
uci set wizard.default.wan_proto="$PROTO"
uci set wizard.default.wan_pppoe_user="$USERNAME"
uci set wizard.default.wan_pppoe_pass="$PASSWORD"

uci commit wizard

exit 0
