# applications/luci-app-apfree-wifidog/root/etc/uci-defaults/luci-app-apfree-wifidog.sh
#!/bin/sh
[ -f "/etc/config/hostnames" ] || {
	echo 'config hostname' > /etc/config/hostnames
}

exit 0
