# luci-app-webd/root/etc/uci-defaults/luci-webd
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@webd[-1]
	add ucitrack webd
	set ucitrack.@webd[-1].init=webd
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
