# small-package/luci-app-LingTiGameAcc/root/etc/uci-defaults/45_luci-lingti
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@lingti[-1]
	add ucitrack lingti
	set ucitrack.@lingti[-1].init=lingti_luci
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
