# small-package/luci-app-tencentddns/files/root/etc/uci-defaults/luci-tencentddns
#!/bin/sh 

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@tencentddns[-1]
	add ucitrack tencentddns
	set ucitrack.@tencentddns[-1].init=tencentddns
	commit ucitrack
EOF

exit 0
