# luci-app-fchomo/root/etc/uci-defaults/99_luci-app-fchomo-migration
#!/bin/sh
# Migration script for fchomo
# Used to migrate LuCI application option.

CONF=fchomo

#uci commit "$CONF"

exit 0
