# mmdvm-host/files/etc/uci-defaults/99_mmdvm-init
#!/bin/sh

[ ! -f /etc/config/mmdvm ] && touch /etc/config/mmdvm
[ -z `uci -q get mmdvm.DMR.Enable` ] && {
	uci set mmdvm.DMR=mmdvmhost
	uci set mmdvm.DMR.Enable=0
	uci commit mmdvm
}

exit 0
