# my-default-settings/target/rockchip/armv8/files/etc/init.d/fa-fancontrol
#!/bin/sh /etc/rc.common

START=21
start() {
    /usr/bin/fa-fancontrol.sh &
}

stop() {
kill -9 $(ps -w | grep fa-fancontrol | grep -v grep | awk '{print$1}') 2>/dev/null
}

restart() {
stop
start
}
