# my-default-settings/files/etc/hotplug.d/iface/95-autoneg
#!/bin/sh
#
for eth in "$(ip address | grep ^[0-9] | awk -F: '{print $2}' | sed "s/ //g" | grep '^[e]' | grep -vE "(@|\.)")"; do
	ethtool -s $eth autoneg on
done
