# luci-app-usb-printer/root/etc/uci-defaults/luci-usb-printer
#!/bin/sh

[ -f /etc/init.d/p910nd ] && /etc/init.d/p910nd disable

exit 0
