# luci-app-parentcontrol/root/etc/hotplug.d/iface/97-parentcontrol
#!/bin/sh
CONFIG=parentcontrol

[ "$ACTION" = ifup ] || exit 0
[ "$(`uci -q get $CONFIG.@basic[0].enabled`)" == 1 ] || exit 0
/etc/init.d/parentcontrol start

