# luci-app-wifi-ap/files/etc/init.d/rc.local
# ...existing code...
/usr/sbin/ap-set-bridge-ip.sh
exit 0
# ...existing code...