# luci-app-bypass/root/etc/uci-defaults/luci-bypass
[ ! -f "/usr/share/ucitrack/luci-app-bypass.json" ] && {
    cat > /usr/share/ucitrack/luci-app-bypass.json << EEOF
{
    "config": "bypass",
    "init": "bypass"
}
EEOF
}
uci batch <<EOF
	delete ucitrack.@bypass[-1]
	add ucitrack bypass
	set ucitrack.@bypass[-1].init=bypass
	commit ucitrack
	delete firewall.bypass
	set firewall.bypass=include
	set firewall.bypass.type=script
	set firewall.bypass.path=/var/etc/bypass.include
	set firewall.bypass.reload=1
	commit firewall
EOF

chmod +x /etc/init.d/bypass /usr/share/bypass/* >/dev/null 2>&1
rm -rf /tmp/luci-*

grep -q ip-api.com /etc/bypass/black.list ||
	sed -i '$a ip-api.com' /etc/bypass/black.list

exit 0
